# rpmconf/Env/crontabs/crontab
# ZIMBRASTART -- DO NOT EDIT ANYTHING BETWEEN THIS LINE AND ZIMBRAEND
#
# Log pruning
#
30 2 * * * find /opt/zimbra/log/ -type f -name \*log\* -mtime +8 -daystart -exec rm {} \; > /dev/null 2>&1
#
# Status logging
#
* * * * * /opt/zimbra/bin/zmstatuslog
*/10 * * * * /opt/zimbra/bin/zmdisklog
