# rpmconf/Env/crontabs/crontab.mta
#
# crontab.mta
#
#
# Queue logging
#
*/10 * * * * /opt/zimbra/bin/zmqueuelog
