# rpmconf/Env/crontabs/crontab
# ZIMBRASTART -- DO NOT EDIT ANYTHING BETWEEN THIS LINE AND ZIMBRAEND
#
# Log pruning
#
30 2 * * * find /opt/zimbra/log/ -type f -name \*log\* -mtime +8 -daystart -exec rm {} \; > /dev/null 2>&1
30 2 * * * find /opt/zimbra/tomcat/logs/ -type f -name \*log\* -mtime +8 -daystart -exec rm {} \; > /dev/null 2>&1
0 2 * * * /opt/zimbra/bin/zmroll_catalina.sh > /dev/null 2>&1
